# deploy/systemd/bblab-tier1.timer
[Unit]
Description=bugbounty-lab Tier 1 timer (diário, fora de pico)

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
