# deploy/systemd/bblab-tier2.timer
[Unit]
Description=bugbounty-lab Tier 2 timer (a cada 2-3 dias, madrugada)

[Timer]
OnCalendar=*-*-1/3 03:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
